# dist/scripts/piscreen-client-localserver.service
[Unit]
Description=piscreen-client Local Webserver
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /home/pi/piscreen-client/player/localserver.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
