# dist/scripts/piscreen-client-networkcheck.service
[Unit]
Description=piscreen-client Python Connection Checker
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /home/pi/piscreen-client/player/checkonline.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
