# dist/files/piscreen-updater.service
[Unit]
Description=PiScreen Updater
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /var/piscreen-client/scripts/update.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
