# dist/files/piscreen.service
[Unit]
Description=PiScreen
After=multi-user.target

[Service]
Type=simple
ExecStart=source /var/piscreen-client/scripts/piscreen.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
