# dist/files/piscreen.service
[Unit]
Description=PiScreen Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /var/piscreen-client/scripts/piscreen.py
Restart=on-abort

[Install]
WantedBy=multi-user.target
